# VR/fixvr/auto/99-valve-index-reboot.rules
# Valve Index HMD: permissions + auto-reboot on plug-in
# Vendor: 28de (Valve) | HMD: 2300

# Permissions: let any user access Valve HID/USB devices (no sudo)
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="28de", MODE="0666", TAG+="uaccess"
SUBSYSTEM=="usb", ATTRS{idVendor}=="28de", MODE="0666", TAG+="uaccess"

# Auto-reboot: send 64-byte HID report (0x16 0x01 + 62 nulls) on HMD plug-in
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="28de", ATTRS{idProduct}=="2300", ACTION=="add", RUN+="/bin/sh -c 'printf \"\\x16\\x01\" | cat - /dev/zero | head -c 64 > /dev/%k'"
